# package/vdi-os/files/etc/systemd/system/vdi-network-setup.service
[Unit]
Description=Bring up network interfaces for VDI installer
After=systemd-udevd.service
Before=vdi-setup-installer.service

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'echo "VDI: Bringing up network interfaces..."; for i in /sys/class/net/{eth*,en*,ib*,wl*}; do [ -e "$i" ] || continue; iface=$(basename $i); echo "VDI: Setting $iface UP"; ip link set $iface up || echo "VDI: Failed to set $iface UP"; done; echo "VDI: Network interfaces status:"; ip link show | grep -E "^[0-9]+:" | head -10'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
